# files/containers.back/redis_1/supervisor.conf.d/redis.opencms.rdb.conf
[program:redis_1]
priority=45
command=redis-server /redis.conf.d/redis.conf
directory=/redis.conf.d/
stopsignal=QUIT
